# build_files/root_files/usr/lib/systemd/system/btrbk.timer.d/override.conf
[Unit]
Description=btrbk hourly snapshots

[Timer]
OnCalendar=hourly
Persistent=false
